# systemd/system-metrics.timer
[Unit]
Description=Timer for System Metrics Collection

[Timer]
OnCalendar=*:*:00
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
